# ansible/roles/authproxy/files/authproxy_supervisor.conf
[program:authproxy]
command=/usr/local/bin/authproxy
redirect_stderr=true
stdout_logfile=/var/log/authproxy.log
